# ansible/roles/promtail/templates/promtail.container.j2
[Unit]
Description=Promtail Log Collector
Documentation=https://grafana.com/docs/loki/latest/clients/promtail/
After=network-online.target
Wants=network-online.target

[Container]
Image={{ promtail_image }}
AutoUpdate=registry

# Network and Ports
Network=host
PublishPort={{ promtail_listen_port }}:{{ promtail_listen_port }}

# Volumes and Configuration
Volume={{ promtail_config_dir }}/promtail-config.yml:/etc/promtail/config.yml:ro
Volume=promtail-data:{{ promtail_data_dir }}:Z
Volume=/var/log:/var/log:ro
Volume=/run/log/journal:/run/log/journal:ro
Volume=/etc/machine-id:/etc/machine-id:ro

# Runtime Configuration
User=root
Group=root

# Command Arguments
Exec=--config.file=/etc/promtail/config.yml

# Podman specific
PodmanArgs=--log-driver=journald
PodmanArgs=--sdnotify=container

# Health Check
HealthCmd=wget --no-verbose --tries=1 --spider http://localhost:{{ promtail_listen_port }}/ready || exit 1
HealthInterval=10s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target