# scripts/backup-timer.service
[Unit]
Description=Homelab Infrastructure Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=william
ExecStart=/home/william/git/podman-homelab/scripts/backup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target